# lucee-nginx/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:lucee]
command=/opt/lucee/tomcat/bin/catalina.sh run
